# crontab
52 * * * * /usr/local/bin/php /motivity/bin/console LoadMultipleFlights --env=prod --no-debug --quiet
